# package/a2dp-app-osmc/files/etc/systemd/system/pulseaudio.service
[Unit]
Description=Pulse Audio
        
[Service]
Type=simple
ExecStart=/usr/bin/pulseaudio
